# deploy/proxmox/carcatcher-watchdog.timer
[Unit]
Description=Run CarCatcher health watchdog every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=15s

[Install]
WantedBy=timers.target
